# subprojects/power-manager/power-cap-action-oem.service
[Unit]
Description= Execute Power Limit OEM exception action

[Service]
ExecStart=echo "Dummy oem exception action"
Type=simple
