# subprojects/nmi-manager/nmi.service
[Unit]
Description= Example to trigger NMI, user has to update this file

[Service]
ExecStart=echo "Example to trigger NMI"
Type=oneshot
