# subprojects/altra-host-error-monitor/xyz.openbmc_project.altra_host_error_monitor.service
[Unit]
Description=Altra Host Error Monitor

[Service]
Restart=always
ExecStart=/usr/sbin/altra-host-error-monitor
Type=simple

[Install]
WantedBy=multi-user.target
